# deploy/systemd/saldoflow-dom.service
[Unit]
Description=SaldoFlow daemon domowy
After=network.target

[Service]
Type=simple
User=marcin
WorkingDirectory=/home/marcin/Pulpit/SaldoFlow
EnvironmentFile=/home/marcin/Pulpit/SaldoFlow/.env
Environment=DAEMON_IMPORT_USER=test
ExecStart=/home/marcin/Pulpit/SaldoFlow/demon/import_watchd /home/marcin/Pulpit/SaldoFlow/runtime/watch/dom /home/marcin/Pulpit/SaldoFlow/.venv/bin/python /home/marcin/Pulpit/SaldoFlow/manage.py test --module dom --status-dir /home/marcin/Pulpit/SaldoFlow/runtime/daemon_status --archive-ok /home/marcin/Pulpit/SaldoFlow/runtime/archive/dom/ok --archive-error /home/marcin/Pulpit/SaldoFlow/runtime/archive/dom/error --foreground
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
